# HOMEDIR/.config/autostart/xfce4-notes-autostart.desktop
[Desktop Entry]
Hidden=false
Exec=xfce4-notes
Terminal=false
